# openclaw/deploy/remitroute-heartbeat.timer
# Fires the RemitRoute money engine every 20 minutes (the product default).
# Install to /etc/systemd/system/ on the VPS, then:
# systemctl enable --now remitroute-heartbeat.timer

[Unit]
Description=RemitRoute heartbeat timer (every 20 minutes)

[Timer]
OnBootSec=2min
OnUnitActiveSec=20min
RandomizedDelaySec=60
Unit=remitroute-heartbeat.service
Persistent=true

[Install]
WantedBy=timers.target
