# openclaw/deploy/remitroute-heartbeat.timer
# Fires the RemitRoute money engine every 2 minutes. Install to
# /etc/systemd/system/ on the VPS, then: systemctl enable --now remitroute-heartbeat.timer

[Unit]
Description=RemitRoute heartbeat timer (every 2 minutes)

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
RandomizedDelaySec=10
Unit=remitroute-heartbeat.service
Persistent=true

[Install]
WantedBy=timers.target
